# utils/reboot.service
[Unit]
Description=Reboot service

[Service]
ExecStart=/usr/bin/reboot

[Install]
WantedBy=multi-user.target